# services.d/nginx/run
#!/command/with-contenv sh

nginx
echo "nginx running"
