# services.d/nginx/run
#!/usr/bin/with-contenv sh

nginx
echo "nginx running"
